# packages/gnome-desktop/gdm/files/systemd/gdm.service
[Unit]
Description=Display Manager
Conflicts=kdm.service slim.service xdm.service

[Service]
ExecStart=/usr/sbin/gdm-binary --no-daemon

[Install]
Alias=displaymanager.service
WantedBy=graphical.target
